# s6/sv/mdevd/run
#!/usr/bin/execlineb -P
exec -c
export PATH /usr/bin
mdevd -D 3
